# .docker/supervisord.conf
; supervisord.conf
; Configuração principal do Supervisor
[supervisord]
nodaemon=true                       ; Importante para o Docker: executa o Supervisor em primeiro plano.
logfile=/var/log/supervisor/supervisord.log ; Caminho para o log principal do Supervisor.
pidfile=/var/run/supervisord.pid    ; Caminho para o arquivo PID do Supervisor.
loglevel=info                       ; Nível de log (debug, info, warn, error, critical).

; Habilita o supervisorctl (socket Unix) para gerenciar/inspecionar processos.
[unix_http_server]
file=/run/supervisord.sock

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; Configuração para o processo PHP-FPM.
; O master do PHP-FPM roda como root (necessário para abrir o error_log em
; /proc/self/fd/2 e gerenciar os workers); os workers caem para www-data
; conforme o pool (user/group em php-fpm.d). Definir user=www-data aqui
; quebrava a inicialização do FPM (Permission denied no error_log).
[program:php-fpm]
command = php-fpm -F              ; Executa o PHP-FPM em primeiro plano (-F).
autostart = true                  ; Inicia automaticamente o PHP-FPM com o Supervisor.
autorestart = true                ; Reinicia o PHP-FPM automaticamente se ele parar.
priority = 5                      ; Prioridade de inicialização (menor é mais cedo).
stdout_logfile=/var/log/supervisor/php-fpm_stdout.log ; Log de saída padrão.
stderr_logfile=/var/log/supervisor/php-fpm_stderr.log ; Log de erro padrão.

; Configuração para o worker de fila do Laravel
; Ajuste o número de processos, tentativas, timeout e sleep conforme sua necessidade.
[program:laravel-queue-worker]
process_name=%(program_name)s_%(process_num)02d ; Nome do processo, com número para múltiplos workers.
command = php /app/artisan queue:work --tries=3 --timeout=60 --sleep=3 ; Comando para o worker.
autostart = true                  ; Inicia automaticamente o worker com o Supervisor.
autorestart = true                ; Reinicia o worker automaticamente se ele parar.
numprocs = 1                      ; Número de processos de worker que serão iniciados.
user = www-data                   ; Usuário sob o qual o worker será executado.
stdout_logfile=/var/log/supervisor/queue-worker_stdout.log ; Log de saída padrão.
stderr_logfile=/var/log/supervisor/queue-worker_stderr.log ; Log de erro padrão.
